# skills/task-orchestrator/templates/task.timer.template
[Unit]
Description=Timer for OpenClaw Task: {{TASK_NAME}}

[Timer]
{{TIMER_DIRECTIVES}}

[Install]
WantedBy=timers.target
