# skills/task-orchestrator/templates/task.timer.template
[Unit]
Description=Timer for OpenClaw Task: {{TASK_NAME}}

[Timer]
OnBootSec=5min
OnUnitActiveSec={{INTERVAL}}
Persistent=true

[Install]
WantedBy=timers.target
