# skills/interrupt-service/interrupt-service.service
[Unit]
Description=OpenClaw Interrupt Service - Unified interrupt manager daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node %h/.openclaw/workspace/skills/interrupt-service/interrupt-service.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=interrupt-service

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
